# debian/yandex-legoa-www.postinst
#!/bin/sh
# postinst script for yandex-legoa-www
#
# see: dh_installdeb(1)

set -e

if [ -e /etc/yandex/environment.type ]; then
    yandex_environment=$(cat /etc/yandex/environment.type)
fi

NAME=legoa

# Setting rights and owners
chown -R www-data:www-data /var/run/yandex/legoa
chown -R www-data:www-data /var/log/yandex/legoa

#DEBHELPER#

NGINX_CONF_FILE=90-${NAME}.conf

nginx-enable-mod ${NGINX_CONF_FILE}
invoke-rc.d nginx reload

exit 0
